# keystone/files/keystone.target
[Unit]
Description=Keystone public service
Requires=keystone-admin.service keystone-public.service
[Install]
WantedBy=multi-user.target
